# soal3/crontab
5 6-23/8 * * 0-5 /bin/bash /home/shift3a.sh
